# rootfs/lib/apk/packages/adb-enablemodem.list
/lib/apk/packages/adb-enablemodem.list
/etc/init.d/adb-enablemodem
