# rootfs/lib/apk/packages/busybox.list
/bin/fgrep
/bin/ln
/bin/dd
/bin/grep
/bin/fsync
/bin/pidof
/bin/rm
/bin/date
/bin/ping
/bin/gzip
/bin/ps
/bin/egrep
/bin/mv
/bin/mpstat
/bin/cp
/bin/mkdir
/bin/chown
/bin/traceroute6
/bin/echo
/bin/base64
/bin/uname
/bin/vi
/bin/sync
/bin/ls
/bin/login
/bin/stat
/bin/usleep
/bin/umount
/bin/zcat
/bin/netstat
/bin/passwd
/bin/busybox
/bin/df
/bin/traceroute
/bin/mktemp
/bin/rmdir
/bin/mknod
/bin/touch
/bin/false
/bin/sh
/bin/ping6
/bin/kill
/bin/nice
/bin/lock
/bin/stty
/bin/dmesg
/bin/gunzip
/bin/ed
/bin/netmsg
/bin/sleep
/bin/chmod
/bin/ash
/bin/sed
/bin/true
/bin/tar
/bin/printenv
/bin/mount
/bin/pwd
/bin/cat
/bin/chgrp
/lib/upgrade/keep.d/busybox
/lib/apk/packages/busybox.rusers
/lib/apk/packages/busybox.list
/usr/bin/uuencode
/usr/bin/less
/usr/bin/truncate
/usr/bin/logger
/usr/bin/uptime
/usr/bin/taskset
/usr/bin/test
/usr/bin/yes
/usr/bin/setsid
/usr/bin/clear
/usr/bin/crontab
/usr/bin/dirname
/usr/bin/groups
/usr/bin/hexdump
/usr/bin/tac
/usr/bin/head
/usr/bin/time
/usr/bin/[
/usr/bin/shuf
/usr/bin/top
/usr/bin/killall
/usr/bin/free
/usr/bin/md5sum
/usr/bin/expr
/usr/bin/realpath
/usr/bin/sha3sum
/usr/bin/dos2unix
/usr/bin/awk
/usr/bin/nc
/usr/bin/wc
/usr/bin/hostid
/usr/bin/logname
/usr/bin/env
/usr/bin/uniq
/usr/bin/which
/usr/bin/readlink
/usr/bin/uudecode
/usr/bin/flock
/usr/bin/find
/usr/bin/seq
/usr/bin/diff
/usr/bin/[[
/usr/bin/sort
/usr/bin/mkfifo
/usr/bin/tee
/usr/bin/pstree
/usr/bin/basename
/usr/bin/pgrep
/usr/bin/sum
/usr/bin/unlink
/usr/bin/patch
/usr/bin/od
/usr/bin/strings
/usr/bin/whoami
/usr/bin/tail
/usr/bin/bunzip2
/usr/bin/cal
/usr/bin/expand
/usr/bin/split
/usr/bin/sha1sum
/usr/bin/timeout
/usr/bin/udhcpc6
/usr/bin/cmp
/usr/bin/install
/usr/bin/comm
/usr/bin/cksum
/usr/bin/printf
/usr/bin/fold
/usr/bin/du
/usr/bin/sha256sum
/usr/bin/reset
/usr/bin/nslookup
/usr/bin/tr
/usr/bin/sha512sum
/usr/bin/unexpand
/usr/bin/bzcat
/usr/bin/tty
/usr/bin/cut
/usr/bin/xargs
/usr/bin/nohup
/usr/bin/telnet
/usr/bin/id
/usr/share/acl.d/ntpd.json
/usr/sbin/crond
/usr/sbin/tftpd
/usr/sbin/brctl
/usr/sbin/arping
/usr/sbin/chroot
/usr/sbin/ntpd
/usr/sbin/ntpd-hotplug
/etc/syslog.conf
/etc/profile.d/busybox-history-file.sh
/etc/init.d/sysntpd
/etc/init.d/cron
/etc/capabilities/ntpd.json
/sbin/pivot_root
/sbin/mkswap
/sbin/ip
/sbin/udhcpc
/sbin/route
/sbin/poweroff
/sbin/swapon
/sbin/halt
/sbin/hwclock
/sbin/switch_root
/sbin/reboot
/sbin/swapoff
/sbin/start-stop-daemon
/sbin/vconfig
/sbin/sysctl
/sbin/ifconfig
/sbin/arp
