# rootfs/lib/apk/packages/dawn.list
/lib/apk/packages/dawn.list
/usr/sbin/dawn
/etc/config/dawn
/etc/init.d/dawn
