# rootfs/lib/apk/packages/haproxy.list
/lib/apk/packages/haproxy.list
/usr/sbin/haproxy
/etc/haproxy.cfg
/etc/init.d/haproxy
