# rootfs/lib/apk/packages/internet-detector.list
/lib/apk/packages/internet-detector.list
/usr/bin/internet-detector
/usr/lib/lua/internet-detector/mod_public_ip.lua
/usr/lib/lua/internet-detector/mod_led_control.lua
/usr/lib/lua/internet-detector/mod_network_restart.lua
/usr/lib/lua/internet-detector/mod_user_scripts.lua
/usr/lib/lua/internet-detector/mod_reboot.lua
/usr/lib/lua/internet-detector/main.lua
/etc/config/internet-detector
/etc/internet-detector/down-script.internet
/etc/internet-detector/public-ip-script.internet
/etc/internet-detector/up-script.internet
/etc/init.d/internet-detector
