# rootfs/lib/apk/packages/luci-app-sms-tool-js.list
/lib/apk/packages/luci-app-sms-tool-js.list
/usr/bin/usr_script1.sh
/usr/share/rpcd/acl.d/luci-app-sms-tool-js.json
/usr/share/luci/menu.d/luci-app-sms-tool-js.json
/etc/modem/ussdcodes.user
/etc/modem/phonebook.user
/etc/modem/atcmmds.user
/etc/config/sms_tool_js
/etc/uci-defaults/off_sms.sh
/etc/uci-defaults/setup_sms_tool_js.sh
/etc/init.d/my_new_sms
/www/luci-static/resources/icons/newdelsms.png
/www/luci-static/resources/view/modem/sendsms.js
/www/luci-static/resources/view/modem/sendat.js
/www/luci-static/resources/view/modem/sms_tool_js.css
/www/luci-static/resources/view/modem/readsms.js
/www/luci-static/resources/view/modem/smsconfig.js
/www/luci-static/resources/view/modem/sendussd.js
/sbin/sms_led.sh
/sbin/new_cron_sync.sh
/sbin/smstool_led.sh
