# rootfs/lib/apk/packages/luci-app-passwall.list
/lib/upgrade/keep.d/luci-app-passwall
/lib/apk/packages/luci-app-passwall.list
/usr/share/ucitrack/luci-app-passwall-server.json
/usr/share/ucitrack/luci-app-passwall.json
/usr/share/rpcd/acl.d/luci-app-passwall.json
/usr/share/passwall/haproxy.lua
/usr/share/passwall/0_default_config
/usr/share/passwall/tasks.sh
/usr/share/passwall/socks_auto_switch.sh
/usr/share/passwall/nftables.sh
/usr/share/passwall/rule_update.lua
/usr/share/passwall/iptables.sh
/usr/share/passwall/subscribe.lua
/usr/share/passwall/test.sh
/usr/share/passwall/rules/proxy_host
/usr/share/passwall/rules/gfwlist
/usr/share/passwall/rules/block_host
/usr/share/passwall/rules/chnroute
/usr/share/passwall/rules/proxy_ip
/usr/share/passwall/rules/direct_ip
/usr/share/passwall/rules/lanlist_ipv6
/usr/share/passwall/rules/domains_excluded
/usr/share/passwall/rules/chnlist
/usr/share/passwall/rules/lanlist_ipv4
/usr/share/passwall/rules/block_ip
/usr/share/passwall/rules/chnroute6
/usr/share/passwall/rules/direct_host
/usr/share/passwall/helper_smartdns_add.lua
/usr/share/passwall/haproxy_check.sh
/usr/share/passwall/helper_dnsmasq.lua
/usr/share/passwall/app.sh
/usr/share/passwall/helper_chinadns_add.lua
/usr/share/passwall/helper_smartdns.sh
/usr/share/passwall/monitor.sh
/usr/lib/lua/luci/controller/passwall.lua
/usr/lib/lua/luci/model/cbi/passwall/server/type/ss-rust.lua
/usr/lib/lua/luci/model/cbi/passwall/server/type/ray.lua
/usr/lib/lua/luci/model/cbi/passwall/server/type/sing-box.lua
/usr/lib/lua/luci/model/cbi/passwall/server/type/hysteria2.lua
/usr/lib/lua/luci/model/cbi/passwall/server/type/ssr.lua
/usr/lib/lua/luci/model/cbi/passwall/server/type/trojan-plus.lua
/usr/lib/lua/luci/model/cbi/passwall/server/type/ss.lua
/usr/lib/lua/luci/model/cbi/passwall/server/type/socks.lua
/usr/lib/lua/luci/model/cbi/passwall/server/user.lua
/usr/lib/lua/luci/model/cbi/passwall/server/index.lua
/usr/lib/lua/luci/model/cbi/passwall/client/haproxy.lua
/usr/lib/lua/luci/model/cbi/passwall/client/rule.lua
/usr/lib/lua/luci/model/cbi/passwall/client/type/ss-rust.lua
/usr/lib/lua/luci/model/cbi/passwall/client/type/naive.lua
/usr/lib/lua/luci/model/cbi/passwall/client/type/ray.lua
/usr/lib/lua/luci/model/cbi/passwall/client/type/sing-box.lua
/usr/lib/lua/luci/model/cbi/passwall/client/type/hysteria2.lua
/usr/lib/lua/luci/model/cbi/passwall/client/type/ssr.lua
/usr/lib/lua/luci/model/cbi/passwall/client/type/trojan-plus.lua
/usr/lib/lua/luci/model/cbi/passwall/client/type/ss.lua
/usr/lib/lua/luci/model/cbi/passwall/client/type/tuic.lua
/usr/lib/lua/luci/model/cbi/passwall/client/other.lua
/usr/lib/lua/luci/model/cbi/passwall/client/app_update.lua
/usr/lib/lua/luci/model/cbi/passwall/client/log.lua
/usr/lib/lua/luci/model/cbi/passwall/client/node_subscribe_config.lua
/usr/lib/lua/luci/model/cbi/passwall/client/node_subscribe.lua
/usr/lib/lua/luci/model/cbi/passwall/client/acl.lua
/usr/lib/lua/luci/model/cbi/passwall/client/node_config.lua
/usr/lib/lua/luci/model/cbi/passwall/client/socks_config.lua
/usr/lib/lua/luci/model/cbi/passwall/client/acl_config.lua
/usr/lib/lua/luci/model/cbi/passwall/client/node_list.lua
/usr/lib/lua/luci/model/cbi/passwall/client/shunt_rules.lua
/usr/lib/lua/luci/model/cbi/passwall/client/rule_list.lua
/usr/lib/lua/luci/model/cbi/passwall/client/global.lua
/usr/lib/lua/luci/view/passwall/node_list/link_share_man.htm
/usr/lib/lua/luci/view/passwall/node_list/link_add_node.htm
/usr/lib/lua/luci/view/passwall/node_list/node_list.htm
/usr/lib/lua/luci/view/passwall/rule/rule_version.htm
/usr/lib/lua/luci/view/passwall/rule_list/js.htm
/usr/lib/lua/luci/view/passwall/socks_auto_switch/btn.htm
/usr/lib/lua/luci/view/passwall/cbi/hidevalue.htm
/usr/lib/lua/luci/view/passwall/server/users_list_status.htm
/usr/lib/lua/luci/view/passwall/server/log.htm
/usr/lib/lua/luci/view/passwall/haproxy/status.htm
/usr/lib/lua/luci/view/passwall/app_update/app_version.htm
/usr/lib/lua/luci/view/passwall/global/footer.htm
/usr/lib/lua/luci/view/passwall/global/faq.htm
/usr/lib/lua/luci/view/passwall/global/status.htm
/usr/lib/lua/luci/view/passwall/global/status_bottom.htm
/usr/lib/lua/luci/view/passwall/global/proxy.htm
/usr/lib/lua/luci/view/passwall/log/backup_restore.htm
/usr/lib/lua/luci/view/passwall/log/log.htm
/usr/lib/lua/luci/view/passwall/acl/view_chinadns_log.htm
/usr/lib/lua/luci/passwall/server_app.lua
/usr/lib/lua/luci/passwall/util_shadowsocks.lua
/usr/lib/lua/luci/passwall/util_naiveproxy.lua
/usr/lib/lua/luci/passwall/util_hysteria2.lua
/usr/lib/lua/luci/passwall/util_trojan.lua
/usr/lib/lua/luci/passwall/util_tuic.lua
/usr/lib/lua/luci/passwall/api.lua
/usr/lib/lua/luci/passwall/util_xray.lua
/usr/lib/lua/luci/passwall/com.lua
/usr/lib/lua/luci/passwall/util_sing-box.lua
/etc/config/passwall_server
/etc/uci-defaults/luci-passwall
/etc/init.d/passwall_server
/etc/init.d/passwall
/etc/hotplug.d/iface/98-passwall
/www/luci-static/resources/qrcode.min.js
/www/luci-static/passwall/img/site_icon1_02.png
/www/luci-static/passwall/img/site_icon1_04.png
/www/luci-static/passwall/img/site_icon_03.png
/www/luci-static/passwall/img/site_icon_04.png
/www/luci-static/passwall/img/site_icon1_01.png
/www/luci-static/passwall/img/site_icon_01.png
/www/luci-static/passwall/img/site_icon1_03.png
/www/luci-static/passwall/img/site_icon_02.png
/www/luci-static/passwall/flags/td.svg
/www/luci-static/passwall/flags/ai.svg
/www/luci-static/passwall/flags/cx.svg
/www/luci-static/passwall/flags/jo.svg
/www/luci-static/passwall/flags/sz.svg
/www/luci-static/passwall/flags/ne.svg
/www/luci-static/passwall/flags/ck.svg
/www/luci-static/passwall/flags/lr.svg
/www/luci-static/passwall/flags/ph.svg
/www/luci-static/passwall/flags/ie.svg
/www/luci-static/passwall/flags/ge.svg
/www/luci-static/passwall/flags/iq.svg
/www/luci-static/passwall/flags/ua.svg
/www/luci-static/passwall/flags/sg.svg
/www/luci-static/passwall/flags/la.svg
/www/luci-static/passwall/flags/ni.svg
/www/luci-static/passwall/flags/eu.svg
/www/luci-static/passwall/flags/gu.svg
/www/luci-static/passwall/flags/nr.svg
/www/luci-static/passwall/flags/is.svg
/www/luci-static/passwall/flags/tr.svg
/www/luci-static/passwall/flags/nc.svg
/www/luci-static/passwall/flags/sm.svg
/www/luci-static/passwall/flags/fr.svg
/www/luci-static/passwall/flags/tn.svg
/www/luci-static/passwall/flags/cz.svg
/www/luci-static/passwall/flags/il.svg
/www/luci-static/passwall/flags/de.svg
/www/luci-static/passwall/flags/lv.svg
/www/luci-static/passwall/flags/km.svg
/www/luci-static/passwall/flags/gd.svg
/www/luci-static/passwall/flags/mo.svg
/www/luci-static/passwall/flags/np.svg
/www/luci-static/passwall/flags/tk.svg
/www/luci-static/passwall/flags/ar.svg
/www/luci-static/passwall/flags/tj.svg
/www/luci-static/passwall/flags/nl.svg
/www/luci-static/passwall/flags/tv.svg
/www/luci-static/passwall/flags/lb.svg
/www/luci-static/passwall/flags/hm.svg
/www/luci-static/passwall/flags/er.svg
/www/luci-static/passwall/flags/vn.svg
/www/luci-static/passwall/flags/gl.svg
/www/luci-static/passwall/flags/nz.svg
/www/luci-static/passwall/flags/fo.svg
/www/luci-static/passwall/flags/gm.svg
/www/luci-static/passwall/flags/jp.svg
/www/luci-static/passwall/flags/gw.svg
/www/luci-static/passwall/flags/mx.svg
/www/luci-static/passwall/flags/kh.svg
/www/luci-static/passwall/flags/au.svg
/www/luci-static/passwall/flags/gq.svg
/www/luci-static/passwall/flags/pe.svg
/www/luci-static/passwall/flags/pg.svg
/www/luci-static/passwall/flags/ma.svg
/www/luci-static/passwall/flags/fi.svg
/www/luci-static/passwall/flags/mk.svg
/www/luci-static/passwall/flags/to.svg
/www/luci-static/passwall/flags/do.svg
/www/luci-static/passwall/flags/gt.svg
/www/luci-static/passwall/flags/so.svg
/www/luci-static/passwall/flags/gr.svg
/www/luci-static/passwall/flags/sy.svg
/www/luci-static/passwall/flags/mz.svg
/www/luci-static/passwall/flags/bf.svg
/www/luci-static/passwall/flags/lu.svg
/www/luci-static/passwall/flags/gn.svg
/www/luci-static/passwall/flags/bh.svg
/www/luci-static/passwall/flags/kn.svg
/www/luci-static/passwall/flags/sd.svg
/www/luci-static/passwall/flags/gb-sct.svg
/www/luci-static/passwall/flags/mv.svg
/www/luci-static/passwall/flags/sl.svg
/www/luci-static/passwall/flags/ae.svg
/www/luci-static/passwall/flags/pr.svg
/www/luci-static/passwall/flags/ws.svg
/www/luci-static/passwall/flags/se.svg
/www/luci-static/passwall/flags/ng.svg
/www/luci-static/passwall/flags/bg.svg
/www/luci-static/passwall/flags/tt.svg
/www/luci-static/passwall/flags/nu.svg
/www/luci-static/passwall/flags/ch.svg
/www/luci-static/passwall/flags/bj.svg
/www/luci-static/passwall/flags/ec.svg
/www/luci-static/passwall/flags/uz.svg
/www/luci-static/passwall/flags/et.svg
/www/luci-static/passwall/flags/uy.svg
/www/luci-static/passwall/flags/tc.svg
/www/luci-static/passwall/flags/sb.svg
/www/luci-static/passwall/flags/gb.svg
/www/luci-static/passwall/flags/ki.svg
/www/luci-static/passwall/flags/ls.svg
/www/luci-static/passwall/flags/mp.svg
/www/luci-static/passwall/flags/rs.svg
/www/luci-static/passwall/flags/je.svg
/www/luci-static/passwall/flags/es-ca.svg
/www/luci-static/passwall/flags/hr.svg
/www/luci-static/passwall/flags/mr.svg
/www/luci-static/passwall/flags/cr.svg
/www/luci-static/passwall/flags/re.svg
/www/luci-static/passwall/flags/bi.svg
/www/luci-static/passwall/flags/dm.svg
/www/luci-static/passwall/flags/mt.svg
/www/luci-static/passwall/flags/kw.svg
/www/luci-static/passwall/flags/ht.svg
/www/luci-static/passwall/flags/loading.svg
/www/luci-static/passwall/flags/mm.svg
/www/luci-static/passwall/flags/dk.svg
/www/luci-static/passwall/flags/rw.svg
/www/luci-static/passwall/flags/ye.svg
/www/luci-static/passwall/flags/it.svg
/www/luci-static/passwall/flags/pk.svg
/www/luci-static/passwall/flags/sj.svg
/www/luci-static/passwall/flags/fk.svg
/www/luci-static/passwall/flags/gy.svg
/www/luci-static/passwall/flags/br.svg
/www/luci-static/passwall/flags/ru.svg
/www/luci-static/passwall/flags/az.svg
/www/luci-static/passwall/flags/pa.svg
/www/luci-static/passwall/flags/pm.svg
/www/luci-static/passwall/flags/be.svg
/www/luci-static/passwall/flags/eg.svg
/www/luci-static/passwall/flags/dz.svg
/www/luci-static/passwall/flags/kz.svg
/www/luci-static/passwall/flags/id.svg
/www/luci-static/passwall/flags/mg.svg
/www/luci-static/passwall/flags/ml.svg
/www/luci-static/passwall/flags/ug.svg
/www/luci-static/passwall/flags/bw.svg
/www/luci-static/passwall/flags/am.svg
/www/luci-static/passwall/flags/co.svg
/www/luci-static/passwall/flags/cd.svg
/www/luci-static/passwall/flags/tl.svg
/www/luci-static/passwall/flags/sr.svg
/www/luci-static/passwall/flags/pf.svg
/www/luci-static/passwall/flags/ax.svg
/www/luci-static/passwall/flags/yt.svg
/www/luci-static/passwall/flags/by.svg
/www/luci-static/passwall/flags/eh.svg
/www/luci-static/passwall/flags/ps.svg
/www/luci-static/passwall/flags/al.svg
/www/luci-static/passwall/flags/tf.svg
/www/luci-static/passwall/flags/mf.svg
/www/luci-static/passwall/flags/kr.svg
/www/luci-static/passwall/flags/af.svg
/www/luci-static/passwall/flags/cl.svg
/www/luci-static/passwall/flags/cm.svg
/www/luci-static/passwall/flags/bd.svg
/www/luci-static/passwall/flags/no.svg
/www/luci-static/passwall/flags/sv.svg
/www/luci-static/passwall/flags/gi.svg
/www/luci-static/passwall/flags/ro.svg
/www/luci-static/passwall/flags/us.svg
/www/luci-static/passwall/flags/es.svg
/www/luci-static/passwall/flags/pl.svg
/www/luci-static/passwall/flags/ao.svg
/www/luci-static/passwall/flags/my.svg
/www/luci-static/passwall/flags/lt.svg
/www/luci-static/passwall/flags/qa.svg
/www/luci-static/passwall/flags/bt.svg
/www/luci-static/passwall/flags/ee.svg
/www/luci-static/passwall/flags/ms.svg
/www/luci-static/passwall/flags/vu.svg
/www/luci-static/passwall/flags/cc.svg
/www/luci-static/passwall/flags/sa.svg
/www/luci-static/passwall/flags/ci.svg
/www/luci-static/passwall/flags/cw.svg
/www/luci-static/passwall/flags/tw.svg
/www/luci-static/passwall/flags/in.svg
/www/luci-static/passwall/flags/ba.svg
/www/luci-static/passwall/flags/wf.svg
/www/luci-static/passwall/flags/gb-eng.svg
/www/luci-static/passwall/flags/ag.svg
/www/luci-static/passwall/flags/pt.svg
/www/luci-static/passwall/flags/bv.svg
/www/luci-static/passwall/flags/ad.svg
/www/luci-static/passwall/flags/mh.svg
/www/luci-static/passwall/flags/gg.svg
/www/luci-static/passwall/flags/ve.svg
/www/luci-static/passwall/flags/bb.svg
/www/luci-static/passwall/flags/cv.svg
/www/luci-static/passwall/flags/gf.svg
/www/luci-static/passwall/flags/sn.svg
/www/luci-static/passwall/flags/gp.svg
/www/luci-static/passwall/flags/si.svg
/www/luci-static/passwall/flags/md.svg
/www/luci-static/passwall/flags/sc.svg
/www/luci-static/passwall/flags/vc.svg
/www/luci-static/passwall/flags/bl.svg
/www/luci-static/passwall/flags/dj.svg
/www/luci-static/passwall/flags/vi.svg
/www/luci-static/passwall/flags/ga.svg
/www/luci-static/passwall/flags/om.svg
/www/luci-static/passwall/flags/as.svg
/www/luci-static/passwall/flags/hu.svg
/www/luci-static/passwall/flags/zw.svg
/www/luci-static/passwall/flags/cf.svg
/www/luci-static/passwall/flags/kg.svg
/www/luci-static/passwall/flags/fm.svg
/www/luci-static/passwall/flags/sk.svg
/www/luci-static/passwall/flags/lc.svg
/www/luci-static/passwall/flags/tm.svg
/www/luci-static/passwall/flags/pw.svg
/www/luci-static/passwall/flags/cn.svg
/www/luci-static/passwall/flags/hk.svg
/www/luci-static/passwall/flags/na.svg
/www/luci-static/passwall/flags/mc.svg
/www/luci-static/passwall/flags/aw.svg
/www/luci-static/passwall/flags/bs.svg
/www/luci-static/passwall/flags/ly.svg
/www/luci-static/passwall/flags/at.svg
/www/luci-static/passwall/flags/hn.svg
/www/luci-static/passwall/flags/ss.svg
/www/luci-static/passwall/flags/tg.svg
/www/luci-static/passwall/flags/fj.svg
/www/luci-static/passwall/flags/gh.svg
/www/luci-static/passwall/flags/bq.svg
/www/luci-static/passwall/flags/ca.svg
/www/luci-static/passwall/flags/st.svg
/www/luci-static/passwall/flags/jm.svg
/www/luci-static/passwall/flags/mq.svg
/www/luci-static/passwall/flags/ke.svg
/www/luci-static/passwall/flags/tz.svg
/www/luci-static/passwall/flags/bn.svg
/www/luci-static/passwall/flags/io.svg
/www/luci-static/passwall/flags/mn.svg
/www/luci-static/passwall/flags/za.svg
/www/luci-static/passwall/flags/kp.svg
/www/luci-static/passwall/flags/cg.svg
/www/luci-static/passwall/flags/th.svg
/www/luci-static/passwall/flags/mu.svg
/www/luci-static/passwall/flags/cy.svg
/www/luci-static/passwall/flags/cu.svg
/www/luci-static/passwall/flags/nf.svg
/www/luci-static/passwall/flags/py.svg
/www/luci-static/passwall/flags/bo.svg
/www/luci-static/passwall/flags/mw.svg
