# rootfs/lib/apk/packages/microsocks.list
/lib/apk/packages/microsocks.list
/usr/bin/microsocks
/etc/config/microsocks
/etc/init.d/microsocks
