# rootfs/lib/apk/packages/modemmanager.list
/lib/apk/packages/modemmanager.list
/lib/netifd/proto/modemmanager.sh
/lib/udev/rules.d/77-mm-fibocom-port-types.rules
/lib/udev/rules.d/77-mm-longcheer-port-types.rules
/lib/udev/rules.d/77-mm-gosuncn-port-types.rules
/lib/udev/rules.d/77-mm-ericsson-mbm.rules
/lib/udev/rules.d/77-mm-dell-port-types.rules
/lib/udev/rules.d/77-mm-sierra.rules
/lib/udev/rules.d/77-mm-broadmobi-port-types.rules
/lib/udev/rules.d/77-mm-ublox-port-types.rules
/lib/udev/rules.d/77-mm-telit-port-types.rules
/lib/udev/rules.d/77-mm-tplink-port-types.rules
/lib/udev/rules.d/77-mm-foxconn-port-types.rules
/lib/udev/rules.d/77-mm-zte-port-types.rules
/lib/udev/rules.d/77-mm-mtk-port-types.rules
/lib/udev/rules.d/77-mm-huawei-net-port-types.rules
/lib/udev/rules.d/77-mm-x22x-port-types.rules
/lib/udev/rules.d/77-mm-qcom-soc.rules
/lib/udev/rules.d/80-mm-candidate.rules
/lib/udev/rules.d/77-mm-quectel-port-types.rules
/lib/udev/rules.d/77-mm-cinterion-port-types.rules
/lib/udev/rules.d/77-mm-linktop-port-types.rules
/lib/udev/rules.d/77-mm-nokia-port-types.rules
/lib/udev/rules.d/77-mm-simtech-port-types.rules
/lib/udev/rules.d/77-mm-haier-port-types.rules
/lib/udev/rules.d/77-mm-dlink-port-types.rules
/usr/bin/mmcli
/usr/share/ModemManager/modemmanager.common
/usr/share/ModemManager/fcc-unlock.available.d/105b
/usr/share/ModemManager/fcc-unlock.available.d/413c:81a8
/usr/share/ModemManager/fcc-unlock.available.d/1199
/usr/share/ModemManager/fcc-unlock.available.d/2c7c:030a
/usr/share/ModemManager/fcc-unlock.available.d/1eac
/usr/share/ModemManager/fcc-unlock.available.d/2c7c
/usr/share/ModemManager/fcc-unlock.available.d/03f0:4e1d
/usr/share/ModemManager/fcc-unlock.available.d/105b:e0ab
/usr/share/ModemManager/fcc-unlock.available.d/1199:9079
/usr/share/ModemManager/fcc-unlock.available.d/1eac:1001
/usr/share/ModemManager/fcc-unlock.available.d/105b:e0c3
/usr/share/ModemManager/fcc-unlock.available.d/413c:81a3
/usr/share/ModemManager/mm-foxconn-t77w968-carrier-mapping.conf
/usr/share/dbus-1/system-services/org.freedesktop.ModemManager1.service
/usr/lib/libmm-glib.so.0.10.0
/usr/lib/libmm-glib.so.0
/usr/lib/ModemManager/connection.d/10-report-down
/usr/sbin/ModemManager-monitor
/usr/sbin/ModemManager-wrapper
/usr/sbin/ModemManager
/etc/init.d/modemmanager
/etc/dbus-1/system.d/org.freedesktop.ModemManager1.conf
/etc/hotplug.d/net/25-modemmanager-net
/etc/hotplug.d/wwan/25-modemmanager-wwan
/etc/hotplug.d/tty/25-modemmanager-tty
