# rootfs/lib/apk/packages/netifd.list
/lib/upgrade/keep.d/netifd
/lib/apk/packages/netifd.list
/lib/network/config.sh
/lib/netifd/proto/dhcp.sh
/lib/netifd/utils.sh
/lib/netifd/netifd-proto.sh
/lib/netifd/dhcp.script
/usr/share/udhcpc/default.script
/usr/libexec/network/packet-steering.uc
/etc/uci-defaults/14_migrate-dhcp-release
/etc/udhcpc.user
/etc/init.d/packet_steering
/etc/init.d/network
/etc/hotplug.d/iface/00-netstate
/sbin/ifup
/sbin/ifstatus
/sbin/netifd
/sbin/devstatus
/sbin/ifdown
