# rootfs/lib/apk/packages/telegrambot.list
/lib/apk/packages/telegrambot.list
/usr/lib/telegrambot/plugins/memory.sh
/usr/lib/telegrambot/plugins/ifconfig.sh
/usr/lib/telegrambot/plugins/wll_list.sh
/usr/lib/telegrambot/plugins/netstat.sh
/usr/lib/telegrambot/plugins/uptime.sh
/usr/lib/telegrambot/plugins/plugins.sh
/usr/lib/telegrambot/plugins/functions/ping.sh
/usr/lib/telegrambot/plugins/functions/get_mac.sh
/usr/lib/telegrambot/plugins/ping.sh
/usr/lib/telegrambot/plugins/wol.sh
/usr/lib/telegrambot/plugins/leases.sh
/usr/lib/telegrambot/plugins/uci.sh
/usr/lib/telegrambot/plugins/kernel.sh
/usr/lib/telegrambot/plugins/swports_list.sh
/usr/lib/telegrambot/plugins/wanip.sh
/usr/lib/telegrambot/plugins/wifi_list.sh
/usr/lib/telegrambot/plugins/opkg.sh
/usr/lib/telegrambot/telegrambot
/etc/config/telegrambot
/etc/init.d/telegrambot
