# rootfs/lib/apk/packages/uhttpd.list
/lib/upgrade/keep.d/uhttpd
/lib/apk/packages/uhttpd.list
/usr/sbin/uhttpd
/etc/config/uhttpd
/etc/init.d/uhttpd
