# rootfs/lib/apk/packages/umdns.list
/lib/apk/packages/umdns.list
/usr/sbin/umdns
/etc/config/umdns
/etc/init.d/umdns
