# rootfs/lib/apk/packages/urllogger.list
/lib/apk/packages/urllogger.list
/usr/sbin/urllogger-ctl
/usr/sbin/urllogger
/etc/init.d/urllogger
