# rootfs/lib/apk/packages/urngd.list
/lib/apk/packages/urngd.list
/etc/init.d/urngd
/sbin/urngd
