# rootfs/lib/apk/packages/uwsgi.list
/lib/apk/packages/uwsgi.list
/usr/sbin/uwsgi
/etc/uwsgi/emperor.ini
/etc/init.d/uwsgi
